# rtk_openwrt/feeds/luci/applications/generic/luci-app-rtk-samba/root/etc/uci-defaults/110_luci-rtk-samba
#!/bin/sh

uci -q batch <<-EOF
	set samba.@samba[0].homes='0'
	set samba.@sambashare[0].users='root,user'	
	commit samba
EOF

exit 0
