# rtk_openwrt/package/realtek/rtkbosa/files/rtkbosa.init
#!/bin/sh /etc/rc.common
START=14

start() {
	echo "Starting rtkbosa ..." > /dev/console
	/bin/sh /etc/rtkbosa.sh
	echo "Ending of rtkbosa ..." > /dev/console
}

