# rtk_openwrt/feeds/luci/modules/generic/luci-mod-rtk-generic/root/etc/uci-defaults/20_luci-mod-rtk-generic
#!/bin/sh

uci -q batch <<-EOT
	set luci.main.lang='en'
        commit luci
EOT

exit 0

